# docker/db-backup/crontab
0 0 * * 0 root /bin/backup.sh >> /var/log/cron.log 2>&1
